# docker/production/worker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; -----------------------------------------------------------------------
; Blizzard API queue – 2 workers. ShouldBeUnique + RateLimited on the
; bnet-api bucket prevents duplicate/runaway fetches; lower concurrency
; trades throughput for CPU/RAM stability.
; -----------------------------------------------------------------------
[program:queue-bnet]
command=php -d memory_limit=256M artisan queue:work redis
    --queue=bnet
    --sleep=3
    --tries=3
    --backoff=60
    --timeout=120
    --memory=256
    --max-jobs=500
    --max-time=3600
directory=/var/www/html
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; -----------------------------------------------------------------------
; Raider.io queue – 2 workers. Rate-limited via rio-api bucket.
; -----------------------------------------------------------------------
[program:queue-rio]
command=php -d memory_limit=256M artisan queue:work redis
    --queue=rio
    --sleep=3
    --tries=3
    --backoff=60
    --timeout=120
    --memory=256
    --max-jobs=500
    --max-time=3600
directory=/var/www/html
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; -----------------------------------------------------------------------
; Warcraft Logs queue – 2 workers. Rate-limited via wcl-api bucket.
; -----------------------------------------------------------------------
[program:queue-wcl]
command=php -d memory_limit=256M artisan queue:work redis
    --queue=wcl
    --sleep=3
    --tries=3
    --backoff=60
    --timeout=120
    --memory=256
    --max-jobs=500
    --max-time=3600
directory=/var/www/html
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; -----------------------------------------------------------------------
; Compile queue – CPU-bound JSON/DTO building. Single worker keeps CPU
; predictable; ShouldBeUnique guarantees one compile per character at a time.
; -----------------------------------------------------------------------
[program:queue-compile]
command=php -d memory_limit=512M artisan queue:work redis
    --queue=compile
    --sleep=3
    --tries=3
    --backoff=30
    --timeout=60
    --memory=512
    --max-jobs=500
    --max-time=3600
directory=/var/www/html
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; -----------------------------------------------------------------------
; AI queue – Gemini analysis. Single worker to avoid simultaneous large
; prompt assemblies; ProcessRaidAnalysisJob is ShouldBeUnique per report.
; --timeout bumped to 1200 to match ProcessRaidAnalysisJob::$timeout.
; -----------------------------------------------------------------------
[program:queue-ai]
command=php -d memory_limit=1024M artisan queue:work redis
    --queue=ai
    --sleep=3
    --tries=2
    --backoff=60
    --timeout=1200
    --memory=1024
    --max-jobs=50
    --max-time=7200
directory=/var/www/html
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; -----------------------------------------------------------------------
; Discord queue – low-latency RSVP interactions; shorter timeout.
; -----------------------------------------------------------------------
[program:queue-discord]
command=php -d memory_limit=256M artisan queue:work redis
    --queue=discord
    --sleep=1
    --tries=3
    --backoff=5
    --timeout=60
    --memory=256
    --max-jobs=500
    --max-time=3600
directory=/var/www/html
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; -----------------------------------------------------------------------
; Default queue – fallback / miscellaneous jobs.
; -----------------------------------------------------------------------
[program:queue-default]
command=php -d memory_limit=256M artisan queue:work redis
    --queue=default
    --sleep=3
    --tries=3
    --backoff=30
    --timeout=120
    --memory=256
    --max-jobs=500
    --max-time=3600
directory=/var/www/html
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
