# docker/production/worker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

; -----------------------------------------------------------------------
; Default queue – character/guild sync, Blizzard API, RaiderIO, WCL jobs
; Two workers for parallelism.
; -----------------------------------------------------------------------
[program:queue-default]
command=php artisan queue:work redis
    --queue=default
    --sleep=3
    --tries=3
    --backoff=30
    --timeout=120
    --max-jobs=500
directory=/var/www/html
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; -----------------------------------------------------------------------
; AI queue – long-running Gemini API analysis jobs.
; Single process to avoid concurrent Gemini requests and rate-limit hits.
; Timeout matches the job's $timeout property (300 s).
; -----------------------------------------------------------------------
[program:queue-ai]
command=php artisan queue:work redis
    --queue=ai
    --sleep=3
    --tries=2
    --backoff=60
    --timeout=300
    --max-jobs=100
directory=/var/www/html
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; -----------------------------------------------------------------------
; Discord queue – low-latency RSVP interactions; shorter timeout.
; Dispatch Discord jobs with ->onQueue('discord') to use this worker.
; -----------------------------------------------------------------------
[program:queue-discord]
command=php artisan queue:work redis
    --queue=discord
    --sleep=1
    --tries=3
    --backoff=5
    --timeout=60
    --max-jobs=500
directory=/var/www/html
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
